# esp.service
[Unit]
Description=ESP32 - tester
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/test/esp_tester_flasher/
ExecStart=/home/test/esp_tester_flasher/start.sh
Restart=always

[Install]
WantedBy=multi-user.target